# startfan.service
[Unit]
Description=Script

[Service]
ExecStart=/usr/bin/fan.sh

[Install]
WantedBy=multi-user.target
